# debian/postinst
#!/bin/sh

set -e

LOGDIR="/var/log/pysieved"

if [ "$1" = "configure" ]; then
  if [ ! -d "$LOGDIR" ]; then
    mkdir -p "$LOGDIR"
    chmod 755 "$LOGDIR"
  fi

  if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload || true
    systemctl enable pysieved.service || true
    systemctl restart pysieved.service || true
  fi
fi

exit 0
